# Gewaechshaus/gewachshaus.service
[Unit]
Description=Gewächshaus Steuerungssystem
After=network.target
StartLimitIntervalSec=60
StartLimitBurst=3

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/gewachshaus
ExecStart=/usr/bin/node /home/pi/gewachshaus/server-simple.js
Restart=always
RestartSec=10
StandardOutput=append:/home/pi/gewachshaus/server.log
StandardError=append:/home/pi/gewachshaus/server.log

# Ressourcenbeschränkungen für Raspberry Pi
MemoryMax=150M
CPUQuota=50%

# Sicherheit
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/pi/gewachshaus

# Umgebungsvariablen
Environment=NODE_ENV=production
Environment=PORT=3001

[Install]
WantedBy=multi-user.target
